# supervisor/service/php-fpm.conf
[program:php7-fpm]
command=php-fpm
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
